# supervisord/supervisord.conf
[supervisord]
loglevel=info
nodaemon=true
pidfile=/var/run/supervisord.pid
user=root

[unix_http_server]
file=/tmp/supervisor.sock 

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:haproxy]
command=haproxy -f /etc/haproxy/haproxy.cfg -db
priority=1000
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:consul-template]
priority=1
autostart=true
autorestart=true
command=/usr/bin/confd -backend stackengine -node 172.31.41.36.8443 -scheme https -auth-token c7d728cb050bb558
#command=/usr/bin/consul-template -consul consul:8500 -template "/etc/consul-template/template.d/haproxy.cfg.template:/etc/haproxy/haproxy.cfg:supervisorctl reload haproxy"
logfile=/dev/stdout
loglevel=info
user=root

